# crontab
05 14   * * 6   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 21   * * 0   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 9    * * 1   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 11   * * 1   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 13   * * 1   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 15   * * 1   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 21   * * 1   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 9    * * 2   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 11   * * 2   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 13   * * 2   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 15   * * 2   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 21   * * 2   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 9    * * 3   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 11   * * 3   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 13   * * 3   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 15   * * 3   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 21   * * 3   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 9    * * 4   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 11   * * 4   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 13   * * 4   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 15   * * 4   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 21   * * 4   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 9    * * 5   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 11   * * 5   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 13   * * 5   ubuntu env DISPLAY=:0 bash $HOME/start.sh
30 15   * * 5   ubuntu env DISPLAY=:0 bash $HOME/start.sh
