# crontab
40 21   * * 0   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 7    * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 9    * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 11   * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 13   * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 15   * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 21   * * 1   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 7    * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 9    * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 11   * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 13   * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 15   * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 21   * * 2   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 7    * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 9    * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 11   * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 13   * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 15   * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 21   * * 3   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 7    * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 9    * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 11   * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 13   * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 15   * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 21   * * 4   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 7    * * 5   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 9    * * 5   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 11   * * 5   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 13   * * 5   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
40 15   * * 5   cd $HOME/TRADE/captura && git pull && /bin/bash $HOME/TRADE/captura/captura.sh
